# script/crontab
RAILS_ENV=production
30 6 * * 1-5 cd /work/railsapps/stock_db_capture && /usr/local/bin/rake active_trader:capture_live_quotes
